# aiops/services/n8n.service
# n8n.service
# autogenerated by Podman 5.4.0
# Sun Mar 16 16:59:39 EDT 2025

[Unit]
Description=Podman n8n.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target n8n-postgres.service
After=network-online.target n8n-postgres.service
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	--replace \
	-d \
	--name n8n \
	-p 5678:5678 \
	--network=pasta:-T,5432,-T,11434 \
	-e DB_TYPE=postgresdb \
	--add-host=postgres.local:127.0.0.1 \
    --add-host=ollama.local:127.0.0.1 \
	-e DB_POSTGRESDB_HOST=postgres.local \
	-e DB_POSTGRESDB_PORT=5432 \
	-e N8N_RUNNERS_ENABLED=true \
	-e TZ=America/New_York \
	-e N8N_ENFORCE_SETTINGS_FILE_PERMISSIONS=true \
	--secret POSTGRES_USER,type=env,target=DB_POSTGRESDB_DATABASE \
	--secret POSTGRES_USER,type=env,target=DB_POSTGRESDB_USER \
	--secret POSTGRES_PASSWORD,type=env,target=DB_POSTGRESDB_PASSWORD \
	-v n8n_data:/home/node/.n8n docker.n8n.io/n8nio/n8n
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
