# aiops/services/qdrant.service
# container-qdrant.service
# autogenerated by Podman 5.4.0
# Mon Mar 31 21:35:18 EDT 2025

[Unit]
Description=Podman container-qdrant.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	--replace \
	-d \
	-p 127.0.0.1:6333:6333 \
	-p 127.0.0.1:6334:6334 \
	--network=pasta:-T,11434,-T,9060,-T,3000 \
	--add-host=ollama.local:127.0.0.1 \
	--name=qdrant \
	-v qdrant:/qdrant/storage docker.io/qdrant/qdrant
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
